# apt-cacher.cron.daily
#!/bin/sh

# Only do anything if apt-cacher is actually installed
if [ ! -x /usr/share/apt-cacher/apt-cacher.pl ]; then
  exit 0
fi

# Run the report generator
/usr/share/apt-cacher/apt-cacher-report.pl >/dev/null

# Run the cache cleaner
/usr/share/apt-cacher/apt-cacher-cleanup.pl 2>/dev/null

