# pkg-cacher.cron.daily
#!/bin/sh

# Only do anything if pkg-cacher is actually installed
if [ ! -x /usr/share/pkg-cacher/pkg-cacher.pl ]; then
  exit 0
fi

# Run the report generator
/usr/share/pkg-cacher/pkg-cacher-report.pl >/dev/null

# Run the cache cleaner
/usr/share/pkg-cacher/pkg-cacher-cleanup.pl 2>/dev/null

